# script10/recipes/core/busybox/files/etc/service/syslog/run
#!/bin/sh
# syslog simples (busybox syslogd)
exec /bin/busybox syslogd -n -O /var/log/messages
