# script9/ports/wayland/wl-clipboard-x11-bridge/files/wl-x11-clipboard-bridge.service
[Unit]
Description=Wayland <-> X11 clipboard bridge (wl-clipboard + xclip)
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/wl-x11-clipboard-bridge
Restart=on-failure

[Install]
WantedBy=default.target
