# script9/ports/base/runit-services/files/etc/sv/seatd/run
#!/bin/sh
exec 2>&1

# Ajuste o grupo conforme seu desktop (muito comum: video).
# Se você criar grupo "seat", troque para: -g seat
exec /usr/bin/seatd -g video
