# script9/ports/base/runit-services/files/etc/sv/dhcpcd/run
#!/bin/sh
exec 2>&1

# dhcpcd normalmente daemoniza; com debug (-d) ele não daemoniza (mantém foreground). 9
# -t 0: espera indefinidamente por lease
exec /usr/bin/dhcpcd -d -t 0
