# script9/ports/base/runit-services/files/etc/sv/elogind/run
#!/bin/sh
exec 2>&1
exec /usr/lib/elogind/elogind 2>/dev/null || exec /usr/sbin/elogind
