# script10/recipes/core/busybox/files/etc/service/getty-tty1/run
#!/bin/sh
# getty no tty1
exec /bin/busybox getty -L 115200 tty1 vt100
